# alari-profile/airootfs/etc/systemd/system/alari-bluetooth-agent.service
[Unit]
Description=ALARI Bluetooth Auto-Pairing Agent
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/local/bin/alari-bluetooth-agent.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
