# alari-profile/airootfs/etc/systemd/system/alari-persist.service
[Unit]
Description=ALARI Persistence Setup
After=local-fs.target
Before=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/alari-persist-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
